# bot.conf
[program:translator]
directory=/home/ubuntu/translator

command=python3 main.python3

autostart=true
autorestart=true
environment=HOME="/home/ubuntu",USER="ubuntu"

stderr_logfile=/home/ubuntu/translator/logfile_err.log
stdout_logfile=/home/ubuntu/translator/logfile.log